# root/etc/services.d/dnsmasq/run
#!/usr/bin/with-contenv sh

exec /usr/sbin/dnsmasq -C /etc/dnsmasq.conf -k;
